# etl/sales-daily/pricing-etl-sales-daily.service
[Unit]
Description=Pricing ETL Sales Daily
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/srv/pricing-data-services
EnvironmentFile=/srv/pricing-data-services/etl/sales-daily/.env
ExecStart=/srv/pricing-data-services/.venv/bin/python -m etl.sales-daily.sales-daily-etl
User=root
Group=root

# Opcional: logs en journal + tu jsonl ya queda en ETL_LOG_PATH
StandardOutput=journal
StandardError=journal